# grommet-gui/supervisord.conf
[supervisord]
nodaemon=true

[program:frontend]
command=service nginx start

[program:backend]
command=python backend.py
